# driver/sl4a-touch.service
[Unit]
Description=Surface Laptop 4 Touchscreen Driver (standard HID with one-shot raw-mode fallback)
After=local-fs.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/lib/sl4a-touch/load-touch.sh
# SAFETY: Never unload kernel modules on a running system (SESSION_HANDOFF.md rule).
# /bin/true is a safe no-op that exits successfully without touching the kernel.
ExecStop=/bin/true

[Install]
WantedBy=multi-user.target
