# driver/sl4a-touch.service
[Unit]
Description=Surface Laptop 4 Touchscreen Driver (spi-amd + spi-hid, standard HID mode)
After=local-fs.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/sbin/modprobe spi-hid
ExecStart=/sbin/modprobe spi-amd
ExecStop=/sbin/modprobe -r spi-hid
ExecStop=/sbin/modprobe -r spi-amd

[Install]
WantedBy=multi-user.target
