# tools/diagnostics/oneshot.service
[Unit]
Description=SL4A Touchscreen One-Shot Test
DefaultDependencies=no
After=local-fs.target
Before=display-manager.service

[Service]
Type=oneshot
ExecStart=/home/simone/spi-hid/tools/diagnostics/oneshot_boot.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
